# psygest-pwa/psygest-pwa.service
[Unit]
Description=PsyGest PWA Patient
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/psygest-pwa
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/pi/psygest-pwa/.env

# Sandboxing léger
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
